# root/etc/s6-overlay/s6-rc.d/init-mod-wg-port-forward-update-conf/up
/etc/s6-overlay/s6-rc.d/init-mod-wg-port-forward-update-conf/run