# root/etc/s6-overlay/s6-rc.d/init-mod-wg-port-forward/run
#!/usr/bin/with-contenv bash

echo "[WG_PORT_FORWARD] looking for wg0.conf"

echo "[WG_PORT_FORWARD] could not find wg0.conf"

echo "[WG_PORT_FORWARD] done"