# root/etc/s6-overlay/s6-rc.d/init-mod-wg-port-forward/run
#!/usr/bin/with-contenv bash

__wg_conf_target="/config/wg0.conf"

echo "[WG_PORT_FORWARD] checking for wg0.conf"

ls "/config"

if [ ! -f "$__wg_conf_target" ]; then
    echo "[WG_PORT_FORWARD] could not find wg0.conf, exiting"
    exit 130
fi

echo "[WG_PORT_FORWARD] reading environment variables"

__post_up_template="iptables -A INPUT -i wg0 -p udp --dport {port_num} -j ACCEPT; iptables -A INPUT -i wg0 -p tcp --dport {port_num} -j ACCEPT; "
__post_up_value=""
__pre_down_template="iptables -D INPUT -i wg0 -p udp --dport {port_num} -j ACCEPT; iptables -D INPUT -i wg0 -p tcp --dport {port_num} -j ACCEPT; "
__pre_down_value=""

while IFS= read -r line; do
  __env_value=${line#*=}
  __env_name=${line%%=*}
  if [[ "$__env_name" =~ ^WG_FORWARDED_PORT_([0-9]+)$ ]]; then
    __forward_number=${BASH_REMATCH[1]}
    __port_number=$__env_value
    echo "[WG_PORT_FORWARD] adding forward number $__forward_number for port $__port_number"
    __post_up_value="$__post_up_value"$( echo $__post_up_template | sed "s/{port_num}/$__port_number/g" )
    __pre_down_value="$__pre_down_value"$( echo $__pre_down_template | sed "s/{port_num}/$__port_number/g" )
  fi
done < <(env)

if [ "$__post_up_value" = "" ]; then
    echo "[WG_PORT_FORWARD] no forwarded ports found, exiting"
    exit 78
fi

echo "[WG_PORT_FORWARD] replacing conf values"

sed "/^PostUp /s/=.*$/= $__post_up_value/;/^PreDown /s/=.*$/= $__pre_down_value/" $__wg_conf_target

cat $__wg_conf_target

echo "[WG_PORT_FORWARD] done"